# container/bridge.conf
[program:envoy-bridge]
command=/usr/bin/node /opt/envoy/ws-bridge.js
directory=/opt/envoy
autostart=true
autorestart=true
priority=100
stdout_logfile=/var/log/neko/envoy-bridge.log
stdout_logfile_maxbytes=10MB
redirect_stderr=true
