# container/websockify.conf
[program:websockify]
command=websockify --web=/usr/share/novnc 8091 localhost:5900
autorestart=true
priority=450
stdout_logfile=/var/log/neko/websockify.log
stdout_logfile_maxbytes=10MB
redirect_stderr=true
